# .config/systemd/user/endor-blog.service
[Unit]
Description=Sync remote endor blog folder to local folder

[Service]
ExecStart=/home/m/.local/bin/update-endor-blog.fish
