# .config/systemd/user/update-rubyrose.service
[Unit]
Description=Sync local files to rubyrose server

[Service]
ExecStart=/home/m/.local/bin/update-rubyrose.fish
